# config/chatbot.service
[Unit]
Description=AI Chatbot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ai_chatbot
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/local/bin/gunicorn -c config/gunicorn_config.py src.api.app_simple:app
Restart=always
RestartSec=10
StandardOutput=append:/root/ai_chatbot/logs/service.log
StandardError=append:/root/ai_chatbot/logs/service_error.log

[Install]
WantedBy=multi-user.target
